# modules/auto_puppet_apply/files/upstart.conf
start on (started dbus and net-device-up IFACE!=lo)
task

console log

exec /opt/rusty-boxen/puppet-apply.sh
