# modules/auto_puppet_apply/files/upstart.conf
start on started dbus

exec /opt/rusty-boxen/puppet-apply.sh
