# modules/rusty_boxen/files/upstart.conf
start on (started dbus and net-device-up IFACE!=lo)
task

console log
chdir /opt/rusty-boxen

exec ./bin/rusty update
